# backup/flowable-backup.timer
[Unit]
Description=Flowable Insight Daily Backup Timer
Requires=flowable-backup.service

[Timer]
# 每天凌晨 2:00 執行
OnCalendar=*-*-* 02:00:00

# 如果錯過執行時間 (例如機器關機)，開機後立即執行
Persistent=true

# 隨機延遲 0-15 分鐘，避免多台機器同時備份
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
